# debian/prerm
#!/bin/sh

systemctl stop crowdsec-blocklist-mirror || echo "cannot stop service"
systemctl disable crowdsec-blocklist-mirror || echo "cannot disable service"
